# integrations/claude-code/mnemo-cc-sync.service.template
[Unit]
Description=Mnemo Cortex sync — push Claude Code session activity to Mnemo
After=network.target

[Service]
Type=simple
ExecStart=%h/github/mnemo-cortex/integrations/claude-code/mnemo-cc-sync-loop.sh
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Override these to match your setup. Leave commented for defaults.
# Environment=MNEMO_URL=http://localhost:50001
# Environment=MNEMO_AGENT_ID=cc
# Environment=MNEMO_CC_SESSIONS_DIR=%h/.claude/projects
# Environment=MNEMO_CC_OFFSET_FILE=%h/.mnemo-cc/cc-sync.offset.json
# Environment=MNEMO_CC_SYNC_INTERVAL=60

[Install]
WantedBy=default.target
